# WindScribe-Upstream/windscribe-watchdog.service
[Unit]
Description=Windscribe Watchdog Check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/windscribe-watchdog.sh
